# laravel-worker_backup.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/globalgatecy/www/services/artisan queue:listen --tries=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=globalgatecy
numprocs=1
redirect_stderr=true
stdout_logfile=/home/globalgatecy/www/services/storage/logs/worker.log
stopwaitsecs=30;

